# board/radxa/rock5b/haos-ota.service

[Unit]
Description=HAOS OTA Update Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/haos-ota-update
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
